# 66-mymodemmy.rules
ACTION=="add", ATTRS{idVendor}=="12d1", ATTRS{idProduct}=="1446", RUN+="/home/Murphy/mymodem/usbmode.sh"